# home/dot_config/autostart/eww-daemon.desktop
[Desktop Entry]
Type=Application
Name=Eww Daemon
Exec=eww daemon
X-GNOME-Autostart-enabled=true
