# home/private_dot_config/autostart/greenclip-daemon.desktop
[Desktop Entry]
Type=Application
Name=Greenclip Daemon
Exec=greenclip daemon
X-GNOME-Autostart-enabled=true
