# home/private_dot_config/autostart/thunar-daemon.desktop
[Desktop Entry]
Type=Application
Name=Thunar Daemon
Exec=thunar --daemon
X-GNOME-Autostart-enabled=true
