# home/private_dot_config/autostart/wal-restore.desktop
[Desktop Entry]
Type=Application
Name=Wal Restore Colorscheme
Exec=wal -R -q
X-GNOME-Autostart-enabled=true
